# raveos-theme/kde/raveos-kde-theme-apply.service
[Unit]
Description=Apply RaveOS KDE theme payload
After=local-fs.target systemd-user-sessions.service
ConditionPathExists=!/var/lib/raveos-kde-theme/.applied

[Service]
Type=oneshot
ExecStart=/usr/lib/raveos-kde-theme/auto-apply.sh
ExecStartPost=/usr/bin/install -Dm644 /dev/null /var/lib/raveos-kde-theme/.applied

[Install]
WantedBy=multi-user.target
