# raveos-hyprland-theme/build/raveos-hyprland-theme/pkg/raveos-hyprland-theme/usr/lib/systemd/system/raveos-hyprland-theme-apply.service
[Unit]
Description=Apply RaveOS Hyprland theme payload
After=local-fs.target systemd-user-sessions.service
ConditionPathExists=!/var/lib/raveos-hyprland-theme/.applied

[Service]
Type=oneshot
ExecStart=/usr/lib/raveos-hyprland-theme/auto-apply.sh
ExecStartPost=/usr/bin/install -Dm644 /dev/null /var/lib/raveos-hyprland-theme/.applied

[Install]
WantedBy=multi-user.target
