# raveos-hyprland-theme/raveos-hyprland-session-init.service
[Unit]
Description=RaveOS Hyprland session init (wallpaper + GTK theme)
After=graphical-session-pre.target
Wants=graphical-session-pre.target
ConditionEnvironment=HYPRLAND_INSTANCE_SIGNATURE

[Service]
Type=oneshot
ExecStart=/usr/lib/raveos-hyprland-theme/session-init.sh
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
